# deploy/supervisor/channel-vault-nas-api.conf
; Channel Vault NAS API under supervisord (bare-metal / VM host install).
; Matches CVN_RESTART_ADAPTER=supervisor with
;   CVN_RESTART_SERVICE_NAME=channel-vault-nas-api
; so the in-app restart runs: supervisorctl restart channel-vault-nas-api
;
; Create the virtualenv first:
;   cd /opt/channel-vault-nas/backend && python3 -m venv .venv && .venv/bin/pip install -e ".[dev]"

[program:channel-vault-nas-api]
command=/opt/channel-vault-nas/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
directory=/opt/channel-vault-nas/backend
user=channelvault
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
; Give in-process yt-dlp passes time to stop cleanly.
stopwaitsecs=30

; Keep the raw API loopback-bound; publish only the web tier via a reverse proxy.
; *_HOST_DIR keys are Docker Compose-only; host installs use the in-process keys
; below. Pin CVN_DATABASE_URL (4 slashes = absolute) so the live DB and backups
; co-locate on your NAS volume instead of defaulting under the install dir.
environment=CVN_DATABASE_URL="sqlite+aiosqlite:////volume1/channel-vault-nas/metadata/app.db",CVN_METADATA_DIR="/volume1/channel-vault-nas/metadata",CVN_DOWNLOAD_DIR="/volume1/channel-vault-nas/archive",CVN_RUNTIME_ENV_FILE="/volume1/channel-vault-nas/runtime/.env.runtime",CVN_AUTH_TOKEN="replace-with-a-long-random-token"

stdout_logfile=/var/log/channel-vault-nas/api.log
redirect_stderr=true
